# cron/crontab
#crontab
@daily web2py *cron/email_incomplete_papers
