# cron/crontab
#crontab
@daily web2py *applications/symposium/cron/email_incomplete_papers.py
